# overlay/etc/cron.d/vtigercrm
* * * * * www-data cd /var/www/vtigercrm; /usr/bin/php -f cron.php > /dev/null 2>&1
